# crontab
# * * * * * /usr/local/bin/python /app/script.py >> /app/script.log 2>&1
* * * * * echo "Morrillo" >> /app/script.log 2>&1
